# topoff/lib/systemd/system/topoff.service
[Unit]
Description=aquarium_topoff
After=network.target

[Service]
ExecStart=/usr/bin/python3 -u topoff.py
WorkingDirectory=/home/lukas/homeauto/
Restart=always
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=topoff
User=lukas

[Install]
WantedBy=multi-user.target
